# files/system/usr/lib/systemd/system/rakuos-pkgroot-sync.service
[Unit]
Description=RakuOS Package Root Sync
Documentation=https://rakuos.org
After=network-online.target systemd-resolved.service rakuos-pkgroot-init.service rakuos-migrate.service
Wants=network-online.target systemd-resolved.service
Requires=rakuos-pkgroot-init.service
ConditionPathExists=/var/lib/rakuos/packages.list

[Service]
Type=oneshot
ExecStartPre=/bin/bash -c 'until nslookup mirrors.fedoraproject.org >/dev/null 2>&1; do echo "Waiting for DNS..."; sleep 2; done'
ExecStart=/usr/libexec/rakuos/rakuos-pkgroot-sync
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
