# files/system/usr/lib/systemd/system/rakuos-updater.service
[Unit]
Description=RakuOS Auto Updater (Topgrade)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/rakuos/rakuos-updater
StandardOutput=journal
StandardError=journal