# files/system/usr/lib/systemd/system/rakuos-pkgroot-services.service
[Unit]
Description=RakuOS Package Root Service Activation
Documentation=https://rakuos.org
After=display-manager.service rakuos-pkgroot-sync.service
Wants=rakuos-pkgroot-sync.service
ConditionPathExists=/var/lib/rakuos/packages.list

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/libexec/rakuos/rakuos-pkgroot-services
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
