# files/system/usr/lib/systemd/system/rakuos-cache-clean.service
[Unit]
Description=RakuOS RPM Cache Cleaner
Documentation=https://rakuos.org
After=rakuos-overlay-sync.service
Wants=rakuos-overlay-sync.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/rakuos/rakuos-cache-clean
StandardOutput=journal
StandardError=journal
