# files/system/usr/lib/systemd/system/rakuos-base-protect.service
[Unit]
Description=RakuOS Base System Protection
Documentation=https://rakuos.org
After=rakuos-overlay-mount.service
Wants=rakuos-overlay-mount.service

[Service]
Type=simple
ExecStart=/usr/libexec/rakuos/rakuos-base-protect
Restart=always
RestartSec=5
TimeoutSec=infinity
KillMode=process
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
